# supporting-files/wpsd-nightly-tasks.service
[Unit]
Description=WPSD Nightly Tasks
Conflicts=wpsd-running-tasks.service

[Service]
ExecStart=/usr/local/sbin/.wpsd-nightly-tasks

