# supporting-files/wpsd-running-tasks.service
[Unit]
Description=WPSD Running Tasks
Conflicts=wpsd-nightly-tasks.service

[Service]
ExecStart=/usr/local/sbin/.wpsd-running-tasks

