# supporting-files/wpsd-nightly-tasks.service
[Unit]
Description=WPSD Nightly Tasks

[Service]
ExecStart=/usr/local/sbin/.wpsd-nightly-tasks
User=root
Group=root
Type=oneshot
