# supporting-files/wpsd-running-tasks.service
[Unit]
Description=WPSD Running Tasks

[Service]
ExecStart=/usr/local/sbin/.wpsd-running-tasks
User=root
Group=root
Type=oneshot
