# supporting-files/wpsd-cache.service
[Unit]
Description=WPSD Sys Cache Service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/.wpsd-sys-cache

