# supporting-files/wpsd-cache.timer
[Unit]
Description=Run WPSD Sys Cache Service Every Minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target

